# config/duan/deploys/dailyrefresh.cron
0 4 * * * /bin/bash /root/dailyrefresh.sh amq
3 4 * * * /bin/bash /root/dailyrefresh.sh exe
5 4 * * * /bin/bash /root/dailyrefresh.sh mwx
