# config/duan/ipspy.cron
*/10 * * * * root bash /opt/duan/ipspy.sh
