# config/duan/ipspy.cron
*/1 * * * * /bin/bash /opt/duan/ipspy.sh
